# deploy/watchdog.service
[Unit]
Description=PiCrawler Watchdog Agent
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/picrawler_project
ExecStart=/usr/bin/python3 /home/pi/picrawler_project/scripts/watchdog.py --metrics-url http://127.0.0.1:8000/metrics --interval 5
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
